# deploy/systemd/aiops.target
[Unit]
Description=KubeNova Platform Target
Wants=aiops-control-api.service aiops-runtime-gateway.service
After=network-online.target

[Install]
WantedBy=multi-user.target
